# bin/gadget-service/gadget.service
[Unit]
Description=Setup Gadget

[Service]
Type=oneshot
ExecStart=gadget-service configure
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
