# lmk-service/hid.service
[Unit]
Description=Setup HID

[Service]
Type=oneshot
ExecStart=/usr/bin/setup-hid.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
