# bin/hid-interface/hid.service
[Unit]
Description=Hid Interface
After=gadget.service

[Service]
Type=exec
ExecStart=hid-interface /lmk/hid
ExecStop=fusermount -u /lmk/hid
Restart=on-failure

[Install]
WantedBy=default.target
